# contrib/init/soteriad.conf
description "Soteria Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env SOTERIAD_BIN="/usr/bin/soteriad"
env SOTERIAD_USER="soteria"
env SOTERIAD_GROUP="soteria"
env SOTERIAD_PIDDIR="/var/run/soteriad"
# upstart can't handle variables constructed with other variables
env SOTERIAD_PIDFILE="/var/run/soteriad/soteriad.pid"
env SOTERIAD_CONFIGFILE="/etc/soteria/soteria.conf"
env SOTERIAD_DATADIR="/var/lib/soteriad"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # soteriad will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$SOTERIAD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run soteriad."
        echo "The setting must appear in $SOTERIAD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Soteria Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$SOTERIAD_PIDDIR"
    chmod 0755 "$SOTERIAD_PIDDIR"
    chown $SOTERIAD_USER:$SOTERIAD_GROUP "$SOTERIAD_PIDDIR"
    chown $SOTERIAD_USER:$SOTERIAD_GROUP "$SOTERIAD_CONFIGFILE"
    chmod 0660 "$SOTERIAD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$SOTERIAD_PIDFILE" \
    --chuid $SOTERIAD_USER:$SOTERIAD_GROUP \
    --exec "$SOTERIAD_BIN" \
    -- \
    -pid="$SOTERIAD_PIDFILE" \
    -conf="$SOTERIAD_CONFIGFILE" \
    -datadir="$SOTERIAD_DATADIR" \
    -disablewallet \
    -daemon

